# services/openrc
#!/sbin/openrc-run

name="PurritoBin"
description="ultra fast, minimalistic command line paste-bin"

command="/usr/bin/purrito"
command_background=true
command_args="${command_args:- -d /var/www/purritobin/}"
pidfile="/run/${RC_SVCNAME}.pid"
command_user="purritobin:purritobin"

depend() {
        need net
}
